# etc/systemd/cal.service
[Unit]
Description=E-Ink Calendar

After=time-sync.target
Wants=time-sync.target

[Service]
ExecStart=/usr/bin/python %h/pytoc/cal.py --log-level INFO --image /dev/shm

[Install]
WantedBy=default.target