# etc/systemd/time-sync.target
[Unit]
Description=System Time Synchronized for user
RefuseManualStart=yes
Requires=basic.target
After=basic.target